# jobs.cron
@hourly /virtualenv/bin/python /path/to/GeoHealthCheck/GeoHealthCheck/healthcheck.py run
@daily /virtualenv/bin/python /path/to/GeoHealthCheck/GeoHealthCheck/models.py flush
