# docker/jobs.cron
@hourly /virtualenv/bin/python /GeoHealthCheck/GeoHealthCheck/models.py run
@daily /virtualenv/bin/python /GeoHealthCheck/GeoHealthCheck/models.py flush

